# scripts/consorsbank-import.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nilsworm.consorsbank-import</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/nilsworm/bin/consorsbank-import.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/nilsworm/Downloads</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/nilsworm/Library/Logs/consorsbank-import.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nilsworm/Library/Logs/consorsbank-import.log</string>
</dict>
</plist>
